# installations/hosts/wombat/sudonters/network/if-up.d/no-naps.sh
#!/bin/sh

if [ "${IFACE}" == "wlan0" ]; then
    /usr/sbin/iw wlan0 set power_save off
fi
